# runtime/ubuntu20.04/assets/components/openssh-server/assets/ssh.service.d/00.options.conf
[Service]
ExecStartPre=rm -f /run/nologin
