# runtime/ubuntu20.04/assets/system/etc/systemd/system/aavm.service
[Unit]
Description=This script setups AAVM

[Service]
Type=oneshot
ExecStart=/bin/bash aavm-init

[Install]
WantedBy=multi-user.target